# whatsweb.service.example
[Unit]
Description=WhatsWeb (Baileys + SQLite + Multiusuarios + Webhook + API)
After=network-online.target
Wants=network-online.target

[Service]
User=root
WorkingDirectory=/opt/whatsweb
ExecStart=/usr/bin/node /opt/whatsweb/server.js
Restart=always
RestartSec=2
StandardOutput=append:/var/log/whatsweb.log
StandardError=append:/var/log/whatsweb.log

# Carregar variáveis de /opt/whatsweb/.env (arquivo NÃO versionado)
EnvironmentFile=/opt/whatsweb/.env

# Alternativa: definir inline (apenas placeholders; evite em produção)
# Environment=NODE_ENV=production
# Environment=TZ=America/Sao_Paulo
# Environment=MASTER_EMAIL=provedor@provedor
# Environment=MASTER_PASSWORD=0123456789
# Environment=RB_TOKEN=COLOQUE_AQUI
# Environment=COOKIE_SECRET=COLOQUE_AQUI_64_HEX

[Install]
WantedBy=multi-user.target
